# supervisord.conf
; supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
; Production process supervisor for the trading bot stack.
;
; Manages three processes:
;   bot      — trading engine (main.py)
;   api      — FastAPI REST + WebSocket server
;   frontend — Vite dev server (swap for 'npm run preview' in production)
;
; Usage:
;   supervisord -c supervisord.conf          # start all in background
;   supervisorctl -c supervisord.conf status # show process states
;   supervisorctl -c supervisord.conf stop bot
;   supervisorctl -c supervisord.conf restart api
;   supervisorctl -c supervisord.conf shutdown
;
; Logs land in ./logs/ — created automatically on first run.
; Set APP_ENV in environment= lines or export before launching.
; ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon         = false
logfile          = logs/supervisord.log
logfile_maxbytes = 10MB
logfile_backups  = 5
loglevel         = info
pidfile          = logs/supervisord.pid
childlogdir      = logs

[unix_http_server]
file = logs/supervisor.sock

[supervisorctl]
serverurl = unix:///%(here)s/logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ─── Bot (trading engine) ─────────────────────────────────────────────────────
[program:bot]
command              = %(here)s/venv/bin/python %(here)s/main.py
directory            = %(here)s
autostart            = true
autorestart          = true
startretries         = 3
startsecs            = 5
stopwaitsecs         = 30
; Graceful shutdown: SIGTERM first, then SIGKILL after stopwaitsecs
stopsignal           = TERM
stdout_logfile       = logs/bot.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups  = 10
stderr_logfile       = logs/bot.err
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups  = 5
; Set APP_ENV here or via the shell before launching supervisord.
; environment = APP_ENV="paper"

; ─── API (FastAPI / uvicorn) ──────────────────────────────────────────────────
[program:api]
command              = %(here)s/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
directory            = %(here)s
autostart            = true
autorestart          = true
startretries         = 3
startsecs            = 3
stopwaitsecs         = 10
stopsignal           = TERM
stdout_logfile       = logs/api.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups  = 10
stderr_logfile       = logs/api.err
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups  = 5

; ─── Frontend (Vite) ─────────────────────────────────────────────────────────
; For development: npm run dev   (hot-reload, port 5173)
; For production:  npm run preview (serves built dist/, port 4173)
; Run 'cd frontend && npm run build' before switching to preview.
[program:frontend]
command              = npm run dev
directory            = %(here)s/frontend
autostart            = true
autorestart          = true
startretries         = 3
startsecs            = 5
stopwaitsecs         = 10
stopsignal           = TERM
stdout_logfile       = logs/frontend.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 5
stderr_logfile       = logs/frontend.err
stderr_logfile_maxbytes = 5MB
stderr_logfile_backups  = 3

; ─── Process group — control all three with one command ──────────────────────
[group:trading]
programs = bot,api,frontend
